# weatherdatalogger/tempest/systemd/tempest-datalogger.service
[Unit]
Description=WeatherFlow Tempest UDP → MQTT Datalogger
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=weatherdatalogger
Group=weatherdatalogger
WorkingDirectory=/opt/weatherdatalogger/tempest
ExecStart=/opt/weatherdatalogger/tempest/venv/bin/python3 tempest_datalogger.py --config /opt/weatherdatalogger/config.ini
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tempest-datalogger

[Install]
WantedBy=multi-user.target
